# systemd_add_lun/src/scsi_scan_target.target
[Unit]
Description=scans given scsi target

[Install]
WantedBy=multi-user.target
